# etc/systemd/system/genius-headers-mirror.service
[Unit]
Description=Genius Overlay - Headers Mirror
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/genius-overlay/.env
WorkingDirectory=/opt/genius-overlay
# Use the compiled JS if you build with tsc:
# ExecStart=/usr/bin/node /opt/genius-overlay/dist/scripts/headers-mirror.js
# Or if you run with tsx:
ExecStart=/usr/bin/node --loader tsx /opt/genius-overlay/scripts/headers-mirror.ts
User=genius
Group=genius
Nice=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
