# etc/systemd/system/gitdata-headers-mirror.service
[Unit]
Description=Gitdata Overlay - Headers Mirror
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/gitdata-overlay/.env
WorkingDirectory=/opt/gitdata-overlay
# Use the compiled JS if you build with tsc:
# ExecStart=/usr/bin/node /opt/gitdata-overlay/dist/scripts/headers-mirror.js
# Or if you run with tsx:
ExecStart=/usr/bin/node --loader tsx /opt/gitdata-overlay/scripts/headers-mirror.ts
User=gitdata
Group=gitdata
Nice=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
